# kafka-cluster.target
[Unit]
Description=Kafka Cluster
Wants=kafka-1.service kafka-2.service kafka-3.service
After=kafka-1.service kafka-2.service kafka-3.service