# ns/ns-coredns.container
[Unit]
Description=Name Server CoreDNS
Requires=ns-dnsproxy.container
After=ns-dnsproxy.container

[Service]
Restart=on-failure

[Container]
Pod=ns.pod
Image=docker.io/coredns/coredns:latest
AutoUpdate=registry
Mount=type=bind,src=/var/containers/ns/coredns/Corefile,dst=/Corefile,relabel=shared
Mount=type=bind,src=/var/containers/ns/coredns/etc/coredns,dst=/etc/coredns,relabel=shared
